# debian/postinst
#!/bin/sh

update-rc.d netcheck defaults || exit 1
invoke-rc.d netcheck start || exit 1
